# etc/motion_notify_example.conf
[program:motion_notify]
command=/path/to/your/venv/bin/python /path/to/the/project/motion_notification_bot.py
user=__YOUR_USR__
directory=/path/to/the/project/
redirect_stderr=true
